# distro/overlays/tablet/all/etc/systemd/system/nftables.service.d/wait-for-conntrack.conf
[Unit]
After=modprobe@nf_conntrack.service
Wants=modprobe@nf_conntrack.service

[Service]
ExecStart=
ExecStart=/usr/local/sbin/nftables-safe
