# distro/overlays/tablet/all/etc/xdg/autostart/maliit-server.desktop
[Desktop Entry]
Type=Application
Name=Maliit Keyboard Server
Exec=maliit-server
OnlyShowIn=Phosh;GNOME;
NoDisplay=true
X-GNOME-AutoRestart=true
