# distro/overlays/tablet/all/etc/xdg/autostart/phosh-lock-on-start.desktop
[Desktop Entry]
Name=Lock Screen on Start
Exec=/usr/local/bin/phosh-lock-on-start
Type=Application
NoDisplay=true
X-GNOME-Autostart-enabled=true
