# distro/overlays/tablet/all/etc/systemd/system/tablet-autologin.service
[Unit]
Description=Enable GDM autologin for gemstone (first boot)
Documentation=https://www.freedesktop.org/software/accountsservice/
After=accounts-daemon.service
Requires=accounts-daemon.service
Before=display-manager.service
ConditionPathExists=!/var/lib/gdm3/.autologin-configured

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/tablet-autologin.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical.target
